# root/etc/uci-defaults/xray_simple
#!/bin/sh
# 1. 确保 xray_simple 默认的单例 general 配置节点存在，若不存在则创建
uci get xray_simple.@general[-1] >/dev/null 2>&1 || uci set xray_simple.general=general

# 2. 确保默认的 xray 资源文件路径 (asset_dir) 已经配置好
uci get xray_simple.@general[0].asset_dir >/dev/null 2>&1 || uci set xray_simple.@general[0].asset_dir='/usr/share/xray'
uci get xray_simple.@general[0].system_log >/dev/null 2>&1 || uci set xray_simple.@general[0].system_log='1'
uci get xray_simple.@general[0].runtime_log_file >/dev/null 2>&1 || uci set xray_simple.@general[0].runtime_log_file='/var/etc/xray-simple/xray.log'
uci get xray_simple.@general[0].dnsmasq_upstream >/dev/null 2>&1 || uci set xray_simple.@general[0].dnsmasq_upstream='0'
uci get xray_simple.@general[0].dnsmasq_xray_port >/dev/null 2>&1 || uci set xray_simple.@general[0].dnsmasq_xray_port='5353'
uci -q get xray_simple.@general[0].lan_ifaces >/dev/null 2>&1 || uci add_list xray_simple.@general[0].lan_ifaces='br-lan'

# 3. 清理已废弃的旧版本临时导入字段，避免残留脏配置
uci -q delete xray_simple.@general[0].import_profile_name
uci -q delete xray_simple.@general[0].import_profile_description
uci -q delete xray_simple.@general[0].import_json

# 4. 清理旧版本曾经添加过的占位符绕过地址，防止重复或冗余规则
uci -q del_list xray_simple.@general[0].bypass_ipv4='0.0.0.0/32'
uci -q del_list xray_simple.@general[0].bypass_ipv6='::/128'

# 5. 提交 xray_simple 配置的变更
uci commit xray_simple

# 6. 不使用 ucitrack 自动 reload。LuCI 保存配置只应写入 UCI；
#    Xray/TProxy 的慢速起停由页面上的显式按钮走后台命令执行。
while uci -q delete ucitrack.@xray_simple[-1]; do :; done
uci -q commit ucitrack
exit 0
