# root/etc/uci-defaults/xray_simple
#!/bin/sh
uci get xray_simple.@general[-1] >/dev/null 2>&1 || uci add xray_simple general >/dev/null 2>&1
uci get xray_simple.@general[0].asset_dir >/dev/null 2>&1 || uci set xray_simple.@general[0].asset_dir='/usr/share/xray'
uci -q del_list xray_simple.@general[0].bypass_ipv4='0.0.0.0/32'
uci -q del_list xray_simple.@general[0].bypass_ipv6='::/128'
uci commit xray_simple
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xray_simple[-1]
	add ucitrack xray_simple
	set ucitrack.@xray_simple[-1].init=xray_simple
	commit ucitrack
EOF
exit 0
