# root/etc/uci-defaults/xray_simple
#!/bin/sh
# 1. 确保 xray_simple 默认的单例 general 配置节点存在，若不存在则创建
uci get xray_simple.@general[-1] >/dev/null 2>&1 || uci add xray_simple general >/dev/null 2>&1

# 2. 确保默认的 xray 资源文件路径 (asset_dir) 已经配置好
uci get xray_simple.@general[0].asset_dir >/dev/null 2>&1 || uci set xray_simple.@general[0].asset_dir='/usr/share/xray'

# 3. 清理已废弃的旧版本临时导入字段，避免残留脏配置
uci -q delete xray_simple.@general[0].import_profile_name
uci -q delete xray_simple.@general[0].import_profile_description
uci -q delete xray_simple.@general[0].import_json

# 4. 清理旧版本曾经添加过的占位符绕过地址，防止重复或冗余规则
uci -q del_list xray_simple.@general[0].bypass_ipv4='0.0.0.0/32'
uci -q del_list xray_simple.@general[0].bypass_ipv6='::/128'

# 5. 提交 xray_simple 配置的变更
uci commit xray_simple

# 6. 将 xray_simple 注册并挂载到系统的 ucitrack 配置，以便在 uci 发生改变时实现对服务的自动状态跟踪与重载触发
uci -q batch <<-EOF >/dev/null
	delete ucitrack.@xray_simple[-1]
	add ucitrack xray_simple
	set ucitrack.@xray_simple[-1].init=xray_simple
	commit ucitrack
EOF
exit 0
