# root/etc/init.d/xray_simple
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=1
NAME=xray_simple

RUNDIR=/var/etc/xray-simple
PIDFILE=/var/run/xray-simple.pid
NFT_FILE="${RUNDIR}/fw4/01_xray_simple.nft"
FW4_INCLUDE=/usr/share/nftables.d/table-pre/xray_simple.nft
ROUTE_LOG="${RUNDIR}/route.log"

uci_get() {
	local ret
	ret="$(uci -q get ${NAME}.@general[0].$1)"
	echo "${ret:-$2}"
}

uci_list() {
	uci -q get ${NAME}.@general[0].$1
}

active_profile() {
	local profile
	profile="$(uci_get active_profile "")"
	if [ -n "$profile" ] && uci -q get "${NAME}.${profile}" >/dev/null 2>&1; then
		echo "$profile"
		return
	fi
	uci -q show "$NAME" | sed -n "s/^${NAME}\.\([^.=]*\)=profile$/\1/p" | head -n 1
}

active_profile_label() {
	local profile name
	profile="$(active_profile)"
	[ -n "$profile" ] || {
		echo "none"
		return
	}
	name="$(profile_get "$profile" name)"
	if [ -n "$name" ]; then
		echo "${name} (${profile})"
	else
		echo "$profile"
	fi
}

profile_get() {
	local section="$1"
	local option="$2"
	uci -q get "${NAME}.${section}.${option}"
}

xray_env() {
	local asset_dir
	asset_dir="$(uci_get asset_dir /usr/share/xray)"
	XRAY_LOCATION_ASSET="$asset_dir" "$@"
}

mark_hex() {
	printf "0x%08x" "$1"
}

nft_set_body() {
	local values="$1"
	local out=""
	local item

	for item in $values; do
		[ -n "$item" ] || continue
		out="${out}${out:+, }${item}"
	done

	echo "$out"
}

nft_string_set_body() {
	local values="$1"
	local out=""
	local item

	for item in $values; do
		[ -n "$item" ] || continue
		out="${out}${out:+, }\"${item}\""
	done

	echo "$out"
}

nft_write_set() {
	local name="$1"
	local type="$2"
	local values="$3"

	cat >> "$NFT_FILE" <<-EOF
	    set ${name} {
	        type ${type}
	        flags interval
	EOF
	if [ -n "$values" ]; then
		cat >> "$NFT_FILE" <<-EOF
	        elements = { ${values} }
	EOF
	fi
	cat >> "$NFT_FILE" <<-EOF
	    }

	EOF
}

write_config() {
	local profile json

	mkdir -p "$RUNDIR"
	profile="$(active_profile)"
	if [ -n "$profile" ]; then
		json="$(profile_get "$profile" json_config)"
	fi
	if [ -z "$json" ]; then
		json="$(uci_get json_config "{}")"
	fi
	printf "%s\n" "$json" > "${RUNDIR}/config.json"
}

write_nft() {
	local tproxy_port mark outbound_mark mark_h outbound_mark_h proxy_lan_dns proxy_router_output
	local bypass_ipv4 bypass_ipv6 nft_mode lan_ifaces bypass_uids bypass_gids wrap_table

	nft_mode="$(uci_get nft_mode firewall4)"
	if [ "$nft_mode" = "direct" ]; then
		NFT_FILE="${RUNDIR}/direct_xray_simple.nft"
		rm -f "${RUNDIR}"/fw4/*.nft
		wrap_table=1
	else
		NFT_FILE="${RUNDIR}/fw4/01_xray_simple.nft"
		rm -f "${RUNDIR}/direct_xray_simple.nft"
		wrap_table=0
	fi
	tproxy_port="$(uci_get tproxy_port 12345)"
	mark="$(uci_get mark 1)"
	outbound_mark="$(uci_get outbound_mark 255)"
	mark_h="$(mark_hex "$mark")"
	outbound_mark_h="$(mark_hex "$outbound_mark")"
	proxy_lan_dns="$(uci_get proxy_lan_dns 1)"
	proxy_router_output="$(uci_get proxy_router_output 1)"
	bypass_ipv4="$(nft_set_body "$(uci_list bypass_ipv4)")"
	bypass_ipv6="$(nft_set_body "$(uci_list bypass_ipv6)")"
	lan_ifaces="$(nft_string_set_body "$(uci_list lan_ifaces)")"
	bypass_uids="$(nft_set_body "$(uci_list bypass_uids)" "")"
	bypass_gids="$(nft_set_body "$(uci_list bypass_gids)" "")"

	mkdir -p "$(dirname "$NFT_FILE")"
	: > "$NFT_FILE"
	if [ "$wrap_table" = "1" ]; then
		cat >> "$NFT_FILE" <<-EOF
	table inet xray_simple {
	EOF
	fi
	nft_write_set "bypass_ipv4" "ipv4_addr" "$bypass_ipv4"
	nft_write_set "bypass_ipv6" "ipv6_addr" "$bypass_ipv6"
	cat >> "$NFT_FILE" <<-EOF
	    chain xray_simple_prerouting {
	        type filter hook prerouting priority mangle +10; policy accept;

	        meta mark ${outbound_mark_h} return
	        ip daddr { 127.0.0.0/8, 224.0.0.0/4, 255.255.255.255/32 } return
	        ip6 daddr { ::1/128, fe80::/10 } return
	EOF

	if [ "$proxy_lan_dns" = "1" ]; then
		if [ -n "$lan_ifaces" ]; then
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol udp iifname { ${lan_ifaces} } udp dport 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr udp iifname { ${lan_ifaces} } udp dport 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		else
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol udp udp dport 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr udp udp dport 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		fi
	fi

	cat >> "$NFT_FILE" <<-EOF
	        ip daddr { 10.0.0.0/8, 172.16.0.0/12, 192.168.0.0/16 } return
	        ip daddr @bypass_ipv4 return
	        ip6 daddr fc00::/7 return
	        ip6 daddr @bypass_ipv6 return
	EOF

	if [ "$proxy_lan_dns" = "1" ]; then
		if [ -n "$lan_ifaces" ]; then
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol tcp iifname { ${lan_ifaces} } meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip protocol udp iifname { ${lan_ifaces} } meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr tcp iifname { ${lan_ifaces} } meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	        ip6 nexthdr udp iifname { ${lan_ifaces} } meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		else
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol tcp meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip protocol udp meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr tcp meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	        ip6 nexthdr udp meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		fi
	else
		if [ -n "$lan_ifaces" ]; then
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol tcp iifname { ${lan_ifaces} } tcp dport != 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip protocol udp iifname { ${lan_ifaces} } udp dport != 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr tcp iifname { ${lan_ifaces} } tcp dport != 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	        ip6 nexthdr udp iifname { ${lan_ifaces} } udp dport != 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		else
			cat >> "$NFT_FILE" <<-EOF
	        ip protocol tcp tcp dport != 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip protocol udp udp dport != 53 meta mark set ${mark_h} tproxy ip to :${tproxy_port} accept
	        ip6 nexthdr tcp tcp dport != 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	        ip6 nexthdr udp udp dport != 53 meta mark set ${mark_h} tproxy ip6 to :${tproxy_port} accept
	EOF
		fi
	fi

	cat >> "$NFT_FILE" <<-EOF
	    }
	EOF

	if [ "$proxy_router_output" = "1" ]; then
		cat >> "$NFT_FILE" <<-EOF

	    chain xray_simple_output {
	        type route hook output priority mangle +10; policy accept;

	        meta mark ${outbound_mark_h} return
	EOF
		if [ -n "$bypass_uids" ]; then
			cat >> "$NFT_FILE" <<-EOF
	        meta skuid { ${bypass_uids} } return
	EOF
		fi
		if [ -n "$bypass_gids" ]; then
			cat >> "$NFT_FILE" <<-EOF
	        meta skgid { ${bypass_gids} } return
	EOF
		fi
		cat >> "$NFT_FILE" <<-EOF
	        ip daddr { 127.0.0.0/8, 224.0.0.0/4, 255.255.255.255/32 } return
	        ip daddr { 10.0.0.0/8, 172.16.0.0/12, 192.168.0.0/16 } return
	        ip daddr @bypass_ipv4 return
	        ip6 daddr { ::1/128, fe80::/10 } return
	        ip6 daddr fc00::/7 return
	        ip6 daddr @bypass_ipv6 return
	        ip protocol tcp meta mark set ${mark_h} accept
	        ip protocol udp meta mark set ${mark_h} accept
	        ip6 nexthdr tcp meta mark set ${mark_h} accept
	        ip6 nexthdr udp meta mark set ${mark_h} accept
	    }
	EOF
	fi

	if [ "$wrap_table" = "1" ]; then
		cat >> "$NFT_FILE" <<-EOF
	}
	EOF
	fi
}

setup_routes() {
	local mark table_v4 table_v6
	mark="$(uci_get mark 1)"
	table_v4="$(uci_get route_table_v4 100)"
	table_v6="$(uci_get route_table_v6 106)"

	mkdir -p "$RUNDIR"
	: > "$ROUTE_LOG"
	flush_routes
	ip rule add fwmark "$mark" table "$table_v4" 2>>"$ROUTE_LOG"
	ip route add local 0.0.0.0/0 dev lo table "$table_v4" 2>>"$ROUTE_LOG"
	ip -6 rule add fwmark "$mark" table "$table_v6" 2>>"$ROUTE_LOG"
	ip -6 route add local ::/0 dev lo table "$table_v6" 2>>"$ROUTE_LOG"
}

flush_routes() {
	local mark table_v4 table_v6
	mark="$(uci_get mark 1)"
	table_v4="$(uci_get route_table_v4 100)"
	table_v6="$(uci_get route_table_v6 106)"

	while ip rule del fwmark "$mark" table "$table_v4" 2>/dev/null; do :; done
	ip route del local 0.0.0.0/0 dev lo table "$table_v4" 2>/dev/null
	while ip -6 rule del fwmark "$mark" table "$table_v6" 2>/dev/null; do :; done
	ip -6 route del local ::/0 dev lo table "$table_v6" 2>/dev/null
}

load_firewall() {
	local nft_mode
	local output
	nft_mode="$(uci_get nft_mode firewall4)"

	if [ "$nft_mode" = "direct" ]; then
		rm -f "${RUNDIR}"/fw4/*.nft
		/etc/init.d/firewall restart >/dev/null 2>&1
		nft delete table inet xray_simple 2>/dev/null
		output="$(nft -f "$NFT_FILE" 2>&1)"
		if [ "$?" != "0" ]; then
			printf "%s\n" "$output"
			return 1
		fi
	else
		nft delete table inet xray_simple 2>/dev/null
		rm -f "${RUNDIR}/direct_xray_simple.nft"
		ensure_fw4_include
		output="$(/etc/init.d/firewall restart 2>&1)"
		if [ "$?" != "0" ]; then
			printf "%s\n" "$output"
			return 1
		fi
	fi
}

flush_firewall() {
	nft delete table inet xray_simple 2>/dev/null
	rm -f "${RUNDIR}/direct_xray_simple.nft"
	rm -f "${RUNDIR}"/fw4/*.nft
	/etc/init.d/firewall restart
}

rollback_firewall() {
	nft delete table inet xray_simple 2>/dev/null
	rm -f "${RUNDIR}/direct_xray_simple.nft"
	rm -f "${RUNDIR}"/fw4/*.nft
	/etc/init.d/firewall restart >/dev/null 2>&1
}

ensure_fw4_include() {
	mkdir -p "$(dirname "$FW4_INCLUDE")"
	printf '%s\n' 'include "/var/etc/xray-simple/fw4/*.nft"' > "$FW4_INCLUDE"
}

start_service() {
	[ "$(uci_get enabled 0)" = "1" ] || return 0

	local xray_bin
	xray_bin="$(uci_get xray_bin /usr/bin/xray)"
	command -v "$xray_bin" >/dev/null 2>&1 || {
		echo "xray binary not found: $xray_bin"
		return 1
	}

	modprobe nft_tproxy 2>/dev/null
	write_config
	test_json_file "$xray_bin" "${RUNDIR}/config.json" || return 1
	write_nft
	setup_routes
	load_firewall || {
		rollback_firewall
		flush_routes
		return 1
	}

	procd_open_instance
	procd_set_param respawn 1 1 0
	procd_set_param command "$xray_bin" run -confdir "$RUNDIR"
	procd_set_param env XRAY_LOCATION_ASSET="$(uci_get asset_dir /usr/share/xray)"
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param pidfile "$PIDFILE"
	procd_set_param file /etc/config/xray_simple
	procd_close_instance
}

stop_service() {
	flush_firewall
	flush_routes
}

start_tproxy() {
	modprobe nft_tproxy 2>/dev/null
	write_nft
	setup_routes
	load_firewall || {
		rollback_firewall
		flush_routes
		return 1
	}
	echo "Xray Simple TProxy rules started."
}

stop_tproxy() {
	flush_firewall
	flush_routes
	echo "Xray Simple TProxy rules stopped."
}

reload_service() {
	stop
	start
}

status_service() {
	local pids table_v6
	pids="$(xray_pids)"
	table_v6="$(uci_get route_table_v6 106)"
	if [ -n "$pids" ]; then
		echo "Xray Simple running, pid ${pids}"
	else
		echo "Xray Simple is not running"
	fi
	echo
	echo "nft mode: $(uci_get nft_mode firewall4)"
	echo "active profile: $(active_profile_label)"
	echo "mark: $(uci_get mark 1)"
	echo "tproxy port: $(uci_get tproxy_port 12345)"
	echo "asset dir: $(uci_get asset_dir /usr/share/xray)"
	echo
	ip rule show | grep -E "fwmark|$(uci_get route_table_v4 100)" || true
	ip -6 rule show | grep -E "fwmark|$(uci_get route_table_v6 106)" || true
	if ! ip -6 route show table "$table_v6" 2>/dev/null | grep -Eq '(^|[[:space:]])local[[:space:]]+(default|::/0)'; then
		echo "warning: IPv6 local route is missing in table ${table_v6}"
	fi
}

xray_pids() {
	ps w 2>/dev/null | awk '$5 ~ /(^|\/)xray$/ && index($0, " run ") && index($0, "confdir /var/etc/xray-simple") { printf "%s%s", sep, $1; sep = " " }'
}

recent_xray_logs() {
	logread 2>/dev/null | grep -iE "xray|xray_simple|xray-simple" | tail -30 || true
}

assert_running() {
	local pids

	sleep 2
	pids="$(xray_pids)"
	if [ -n "$pids" ]; then
		echo "Xray Simple running, pid ${pids}"
		return 0
	fi

	echo "Xray Simple start failed: process is not running after start request."
	echo
	echo "Recent Xray logs:"
	recent_xray_logs
	return 1
}

test_config() {
	local xray_bin
	local output
	xray_bin="$(uci_get xray_bin /usr/bin/xray)"
	command -v "$xray_bin" >/dev/null 2>&1 || {
		echo "xray binary not found: $xray_bin"
		return 1
	}
	write_config
	output="$(xray_env "$xray_bin" test -confdir "$RUNDIR" 2>&1)"
	if [ "$?" = "0" ]; then
		printf "%s\n" "$output"
		return 0
	fi

	output="$(xray_env "$xray_bin" run -test -confdir "$RUNDIR" 2>&1)"
	if [ "$?" = "0" ]; then
		printf "%s\n" "$output"
		return 0
	fi

	printf "%s\n" "$output"
	return 1
}

test_json_file() {
	local xray_bin="$1"
	local json_file="$2"
	local test_dir output

	if [ -z "$json_file" ]; then
		json_file="$xray_bin"
		xray_bin="$(uci_get xray_bin /usr/bin/xray)"
	fi
	command -v "$xray_bin" >/dev/null 2>&1 || {
		echo "xray binary not found: $xray_bin"
		return 1
	}
	[ -s "$json_file" ] || {
		echo "JSON file is empty or missing: $json_file"
		return 1
	}

	test_dir="/tmp/xray-simple-test.$$"
	mkdir -p "$test_dir" || return 1
	cp "$json_file" "${test_dir}/config.json" || {
		rm -rf "$test_dir"
		return 1
	}

	output="$(xray_env "$xray_bin" test -confdir "$test_dir" 2>&1)"
	if [ "$?" = "0" ]; then
		rm -rf "$test_dir"
		printf "%s\n" "$output"
		return 0
	fi

	output="$(xray_env "$xray_bin" run -test -confdir "$test_dir" 2>&1)"
	if [ "$?" = "0" ]; then
		rm -rf "$test_dir"
		printf "%s\n" "$output"
		return 0
	fi

	rm -rf "$test_dir"
	printf "%s\n" "$output"
	return 1
}

show_nft_chain() {
	local family="$1"
	local table="$2"
	local chain="$3"
	local output

	output="$(nft list chain "$family" "$table" "$chain" 2>&1)"
	if [ "$?" = "0" ]; then
		printf "%s\n" "$output"
	else
		echo "${family} ${table} ${chain}: not loaded"
	fi
}

show_nft_table() {
	local family="$1"
	local table="$2"
	local output

	output="$(nft list table "$family" "$table" 2>&1)"
	if [ "$?" = "0" ]; then
		printf "%s\n" "$output"
	else
		echo "${family} ${table}: not loaded"
	fi
}

nft_status() {
	local nft_mode
	nft_mode="$(uci_get nft_mode firewall4)"

	echo "=== nft mode ==="
	echo "$nft_mode"
	echo
	if [ "$nft_mode" = "direct" ]; then
		echo "=== direct nft table inet xray_simple ==="
		show_nft_table inet xray_simple
		echo
		echo "=== firewall4 xray_simple chains ==="
		show_nft_chain inet fw4 xray_simple_prerouting
		show_nft_chain inet fw4 xray_simple_output
		echo
		echo "=== generated direct rules ==="
		cat "${RUNDIR}/direct_xray_simple.nft" 2>/dev/null || echo "not generated"
	else
		echo "=== direct nft table inet xray_simple ==="
		show_nft_table inet xray_simple
		echo
		echo "=== firewall4 xray_simple chains ==="
		show_nft_chain inet fw4 xray_simple_prerouting
		show_nft_chain inet fw4 xray_simple_output
		echo
		echo "=== firewall4 include hook ==="
		cat "$FW4_INCLUDE" 2>/dev/null || echo "missing: $FW4_INCLUDE"
		echo
		echo "=== firewall4 table chain summary ==="
		nft list table inet fw4 2>&1 | grep -E "chain xray_simple|hook prerouting|hook output|include" || true
		echo
		echo "=== generated firewall4 rules ==="
		cat "${RUNDIR}/fw4/01_xray_simple.nft" 2>/dev/null || echo "not generated"
	fi
	echo
	echo "=== IPv4 rules/routes ==="
	ip rule show 2>&1 | grep -E "fwmark|$(uci_get route_table_v4 100)" || true
	ip route show table "$(uci_get route_table_v4 100)" 2>&1 || true
	echo
	echo "=== IPv6 rules/routes ==="
	ip -6 rule show 2>&1 | grep -E "fwmark|$(uci_get route_table_v6 106)" || true
	ip -6 route show table "$(uci_get route_table_v6 106)" 2>&1 || true
	if ! ip -6 route show table "$(uci_get route_table_v6 106)" 2>/dev/null | grep -Eq '(^|[[:space:]])local[[:space:]]+(default|::/0)'; then
		echo "warning: IPv6 local route is missing in table $(uci_get route_table_v6 106)"
	fi
	if [ -s "$ROUTE_LOG" ]; then
		echo
		echo "=== route setup log ==="
		cat "$ROUTE_LOG"
	fi
}

geodata_status() {
	local asset_dir
	asset_dir="$(uci_get asset_dir /usr/share/xray)"
	echo "asset_dir: ${asset_dir}"
	if [ -s "${asset_dir}/geoip.dat" ]; then
		echo "geoip: found"
	else
		echo "geoip: missing"
	fi
	if [ -s "${asset_dir}/geosite.dat" ]; then
		echo "geosite: found"
	else
		echo "geosite: missing"
	fi
}

switch_profile() {
	local profile="$1"
	[ -n "$profile" ] || {
		echo "missing profile section id"
		return 1
	}
	uci -q get "${NAME}.${profile}" >/dev/null || {
		echo "profile not found: $profile"
		return 1
	}
	uci set "${NAME}.@general[0].active_profile=${profile}"
	uci set "${NAME}.@general[0].enabled=1"
	uci commit "$NAME"
	reload_service || return 1
	assert_running || return 1
	echo "active profile switched to $(active_profile_label)"
}

import_profile() {
	local profile_name="$1"
	local profile_description="$2"
	local json_file="${3:-/tmp/xray-simple-import.json}"
	local profile json output

	case "$json_file" in
		/tmp/xray-simple-import.json) ;;
		*)
			echo "unsupported import path: $json_file"
			return 1
			;;
	esac

	profile_name="$(printf "%s" "$profile_name" | tr '\r\n' '  ')"
	profile_description="$(printf "%s" "$profile_description" | tr '\r\n' '  ')"
	[ -n "$profile_name" ] || profile_name="imported"

	output="$(test_json_file "$json_file" 2>&1)" || {
		printf "%s\n" "$output"
		return 1
	}
	json="$(cat "$json_file")" || return 1

	profile="$(uci add "$NAME" profile)" || {
		echo "failed to create profile"
		return 1
	}
	uci set "${NAME}.${profile}.name=${profile_name}"
	if [ -n "$profile_description" ]; then
		uci set "${NAME}.${profile}.description=${profile_description}"
	else
		uci -q delete "${NAME}.${profile}.description"
	fi
	uci set "${NAME}.${profile}.json_config=${json}"
	uci set "${NAME}.@general[0].active_profile=${profile}"
	uci commit "$NAME" || return 1

	echo "imported profile: ${profile_name} (${profile})"
}

export_profile() {
	local profile="$1"
	local json
	[ -n "$profile" ] || profile="$(active_profile)"
	if [ -n "$profile" ]; then
		json="$(profile_get "$profile" json_config)"
	fi
	if [ -z "$json" ]; then
		json="$(uci_get json_config "{}")"
	fi
	printf "%s\n" "$json"
}

start_now() {
	uci set "${NAME}.@general[0].enabled=1"
	uci commit "$NAME"
	if [ -n "$(xray_pids)" ]; then
		reload_service || return 1
	else
		start || return 1
	fi
	assert_running || return 1
	echo "Xray Simple start requested."
}

stop_now() {
	uci set "${NAME}.@general[0].enabled=0"
	uci commit "$NAME"
	stop
	echo "Xray Simple stopped and disabled in Xray Simple settings."
}

restart_now() {
	uci set "${NAME}.@general[0].enabled=1"
	uci commit "$NAME"
	reload_service || return 1
	assert_running || return 1
	echo "Xray Simple restart requested."
}

extra_command "test_config" "Validate active Xray JSON configuration"
extra_command "test_json_file" "Validate a JSON config file"
extra_command "nft_status" "Show nftables and policy routing status"
extra_command "geodata_status" "Show Xray geo database status"
extra_command "switch_profile" "Switch active profile and reload"
extra_command "import_profile" "Import a validated profile JSON file"
extra_command "export_profile" "Print active or selected profile JSON"
extra_command "start_now" "Enable in UCI and start"
extra_command "stop_now" "Stop and disable in UCI"
extra_command "restart_now" "Enable in UCI and restart"
extra_command "start_tproxy" "Start TProxy nftables and policy routing only"
extra_command "stop_tproxy" "Stop TProxy nftables and policy routing only"

service_triggers() {
	procd_add_reload_trigger "xray_simple"
}
